# crontab
# m h  dom mon dow   command
00 20 * * * /home/sylecn/bin/mbackup -q

